# packages/kautolog/kautolog-1.1.13.tar.gz/kautolog-1.1.13/src/kautolog/templates/systemd/autologger-rclone-sync.service
[Unit]
Description=Sync terminal logs with rclone

[Service]
Type=oneshot
ExecStart=/usr/bin/rclone sync --fast-list --transfers=4 --checkers=8 --create-empty-src-dirs __KAUTOLOG_ABS_LOGDIR__ __KAUTOLOG_REMOTE__
