# packages/kautolog/kautolog-1.1.13-py3-none-any.whl/kautolog/templates/systemd/autologger-rclone-sync.timer
[Unit]
Description=Run autologger rclone sync periodically

[Timer]
OnBootSec=2m
OnUnitActiveSec=__KAUTOLOG_INTERVAL__min
Persistent=true

[Install]
WantedBy=default.target
